# data/discord/webhooktest/maincron.cron
45 5 * * 1-5 cd /opt/app && node build/scripts/execJedalnicekPost.js